# systemd/tv-automator-xhost.service
[Unit]
Description=Grant Docker containers access to the local X display (for TV-Automator)
# Start after the display manager has fully started its X session
After=graphical.target display-manager.service
Wants=graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=nrm

# Wait up to 30 s for Xwayland to be ready.
# On Wayland+Xwayland (GNOME), there is no /tmp/.X11-unix/X0 file —
# the display is an abstract socket. Wait for the mutter auth file instead.
ExecStartPre=/bin/bash -c '\
  for i in $(seq 1 30); do \
    ls /run/user/1000/.mutter-Xwaylandauth.* 2>/dev/null | grep -q . && exit 0; \
    [ -S /tmp/.X11-unix/X0 ] && exit 0; \
    sleep 1; \
  done; \
  echo "Timed out waiting for Xwayland (mutter auth or X0 socket)"; exit 1'

# Find the right Xauthority and grant access.
# On Wayland+Xwayland (GNOME default), merge the mutter session auth
# into ~/.Xauthority so Docker has a stable path to mount.
ExecStart=/bin/bash -c '\
  export DISPLAY=:0; \
  MUTTER=$(ls /run/user/1000/.mutter-Xwaylandauth.* 2>/dev/null | head -1); \
  if [ -n "$$MUTTER" ] && [ -f "$$MUTTER" ]; then \
    echo "Merging Wayland/Xwayland auth: $$MUTTER"; \
    xauth merge "$$MUTTER"; \
    export XAUTHORITY=/home/nrm/.Xauthority; \
  else \
    for f in \
      /run/user/1000/gdm/Xauthority \
      /run/user/1000/Xauthority \
      /home/nrm/.Xauthority \
      /var/lib/gdm3/.Xauthority \
      /var/lib/gdm/.Xauthority \
      /var/run/lightdm/root/:0 \
      /tmp/.Xauthority; \
    do \
      if [ -f "$$f" ]; then export XAUTHORITY="$$f"; break; fi; \
    done; \
  fi; \
  if [ -z "$$XAUTHORITY" ]; then \
    echo "No Xauthority file found — cannot grant Docker display access"; \
    exit 1; \
  fi; \
  echo "Using XAUTHORITY=$$XAUTHORITY"; \
  xhost +local: && echo "xhost: Docker granted X access on $$DISPLAY"'

# Revoke on stop (optional — comment out if you want access to persist after stop)
ExecStop=/bin/bash -c '\
  export DISPLAY=:0; \
  for f in \
    /run/user/1000/gdm/Xauthority \
    /run/user/1000/Xauthority \
    /home/nrm/.Xauthority \
    /var/lib/gdm3/.Xauthority; \
  do \
    [ -f "$$f" ] && export XAUTHORITY="$$f" && break; \
  done; \
  xhost -local: 2>/dev/null || true'

[Install]
WantedBy=graphical.target
